# systemd/.config/systemd/user/default.target.wants/dunst.service
[Unit]
Description=display dope-ass messages
[Service]
ExecStart=/usr/bin/dunst
[Install]
WantedBy=default.target

